# templates/systemd.service.erb
[Unit]
Description=devpi
Wants=basic.target
After=basic.target network.target

[Service]
User=<%= scope['devpi::user'] %>
Group=<%= scope['devpi::group'] %>
ExecStart=<%= scope['devpi::exec_path'] %> -c  <%= scope['devpi::config_file'] %> --ldap-config  <%= scope['devpi::config_file'] %>
<% if scope['devpi::proxy'] -%>
Environment=http_proxy=<%= scope['devpi::proxy'] %>
Environment=https_proxy=<%= scope['devpi::proxy'] %>
<% end -%>

[Install]
WantedBy=multi-user.target
