# templates/systemd.service.erb
[Unit]
Description=devpi
Wants=basic.target
After=basic.target network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
#ExecStart=<%= @devpi_path %> --host <%= @listen_host %> --port <%= @listen_port %> --serverdir <%= @server_dir %> --config_file 
ExecStart=<%= @devpi_path %> --config_file <%= @devpi_config %>
<% if @proxy -%>
Environment=http_proxy=<%= @proxy %>
Environment=https_proxy=<%= @proxy %>
<% end -%>

[Install]
WantedBy=multi-user.target
