# templates/systemd.service.erb
[Unit]
Description=devpi
Wants=basic.target
After=basic.target network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
ExecStart=<%= @devpi_path %> -c  <%= @devpi_config %> --ldap-config  <%= @devpi_ldap_config %>
<% if @proxy -%>
Environment=http_proxy=<%= @proxy %>
Environment=https_proxy=<%= @proxy %>
<% end -%>

[Install]
WantedBy=multi-user.target
